# di/profile.d/more-historysize.sh
export HISTSIZE=999999999

#for bash
export HISTFILESIZE="${HISTSIZE}"

#for zsh
export SAVEHIST="${HISTFILESIZE}"
